# srcs/influxdb/srcs/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:influx]
command=sh -c "influxd && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command=sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"