# srcs/mysql/srcs/supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:mysqld]
command = sh -c "mysqld -u root && kill -s SIGTERM $(cat supervisord.pid)"
